# config/systemd-user/dotfiles-sync.service
[Unit]
Description=Dotfiles: config sync (commit + push) at login
# Best effort: if the network is not ready yet, only the push fails.
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# git reads its config natively from ~/.config/git/config (SSH signing etc.).
ExecStart=%h/.local/bin/dotfiles_sync

[Install]
WantedBy=default.target
