# config/systemd-user/dotfiles-sync.service
[Unit]
Description=Dotfiles: config sync (commit + push) at login
# Best effort: if the network is not ready yet, only the push fails.
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Same git config as in the interactive shell (SSH signing etc.).
Environment=GIT_CONFIG_GLOBAL=%h/.config/git/.gitconfig
ExecStart=%h/.local/bin/dotfiles_sync

[Install]
WantedBy=default.target
