# config/systemd-system/getty@tty1.service.d/autologin.conf
# autologin user with the name of _leo_, change it to match your user

[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin leo --noclear %I $TERM
