# config/systemd-system/getty@tty1.service.d/autologin.conf
[Service]
# Autologin `leo` on tty1 (no display manager). Clearing ExecStart first is
# required because a drop-in can only *append* ExecStart lines - the empty
# assignment resets the inherited getty@.service command before the override.
# agetty logs the user straight in without a password prompt; the login shell
# then sources ~/.bash_profile, which execs the dwl session on tty1 (see
# config/bash/.bash_profile). `--noclear` keeps the boot messages on screen.
ExecStart=
ExecStart=-/usr/bin/agetty --autologin leo --noclear %I $TERM
