# systemd-user/.config/systemd/user/battery-check.service
[Unit]
Description=Battery level check

[Service]
Type=oneshot
ExecStart=%h/.local/bin/bat_check

[Install]
WantedBy=default.target
