# config/systemd-system/ly@tty2.service.d/keymap.conf
[Service]
# Colemak-DH in the ly password field / on the text VTs.
#
# Boot race: systemd-vconsole-setup runs at boot (~28 s) while the initramfs
# `kms` hook is initializing the console via a modeset. The VT is not yet in
# KD_TEXT/K_XLATE, so vconsole-setup skips setting the keymap entirely
# ("Configuration of first virtual console was skipped, ignoring remaining
# ones") — and the modeset additionally resets the keymap loaded in the initramfs
# to the kernel default (us). Result: the ly login field is on QWERTY instead of
# mod-dh-iso-uk. (A *manual* `systemctl restart systemd-vconsole-setup` after boot
# works because the modeset is done by then — hence the seemingly contradictory
# diagnosis.)
#
# ly itself does not touch the keymap (no xkb_/keymap option, no keymap strings
# in the binary) — it only mirrors the global console keymap. So we load it again
# right before ly-dm: ly starts ~29 s, i.e. after the modeset, where loadkeys
# takes effect reliably again. KEYMAP must match /etc/vconsole.conf
# (config/vconsole/vconsole.conf).
#
# IMPORTANT: like wait-home.conf, a REAL copy on the root partition, NOT
# symlinked into /home via links.conf (systemd reads unit drop-ins at early
# manager start, a /home symlink would still be dead there). Deploy: see README
# "Manual system state".
ExecStartPre=/usr/bin/loadkeys mod-dh-iso-uk
