# config/systemd-system/ly@tty2.service.d/keymap.conf
[Service]
# Colemak-DH im ly-Passwortfeld / auf den Text-VTs.
#
# Boot-Race: systemd-vconsole-setup läuft beim Boot (~28 s), während der
# initramfs-`kms`-Hook die Konsole gerade per Modeset initialisiert. Die VT ist
# dann noch nicht in KD_TEXT/K_XLATE, also überspringt vconsole-setup das Setzen
# der Keymap komplett ("Configuration of first virtual console was skipped,
# ignoring remaining ones") — und der Modeset resettet zusätzlich die im
# initramfs geladene Keymap auf den Kernel-Default (us). Ergebnis: das
# ly-Login-Feld liegt auf QWERTY statt mod-dh-iso-uk. (Ein *manueller*
# `systemctl restart systemd-vconsole-setup` nach dem Boot funktioniert, weil
# der Modeset dann fertig ist — daher die scheinbar widersprüchliche Diagnose.)
#
# ly selbst fasst die Keymap nicht an (keine xkb_/keymap-Option, keine
# keymap-Strings im Binary) — es spiegelt nur die globale Konsolen-Keymap.
# Deshalb laden wir sie unmittelbar vor ly-dm erneut: ly startet ~29 s, also
# nach dem Modeset, da greift loadkeys wieder zuverlässig. KEYMAP muss zu
# /etc/vconsole.conf passen (config/vconsole/vconsole.conf).
#
# WICHTIG: wie wait-home.conf eine REALE Kopie auf der Root-Partition, NICHT per
# links.conf ins /home symlinken (systemd liest Unit-Drop-ins beim frühen
# Manager-Start, ein /home-Symlink wäre da noch tot). Deploy: siehe README
# »Manual system state«.
ExecStartPre=/usr/bin/loadkeys mod-dh-iso-uk
