# config/systemd-system/ly@tty2.service.d/wait-home.conf
[Unit]
# /etc/ly/config.ini is a symlink into the separate /home (repo under
# /home/<user>/dotfiles). ly@tty2 has no built-in dependency on the /home mount
# and starts at boot almost at the same time as it - if it fires too early it
# reads the dead symlink config and falls back to the ly defaults (including
# session_log in $HOME, wrong looks/options). RequiresMountsFor pulls in
# home.mount as an ordering AND requires dependency, so ly only starts after the
# mount.
#
# IMPORTANT: this drop-in file must sit REAL on the root partition (not symlinked
# into /home via links.conf) - systemd reads unit drop-ins at early manager
# start, where a /home symlink would be just as dead as the config itself. So it
# is not in links.conf; it is deployed as a real copy (see README "Manual system
# state").
RequiresMountsFor=/home
