# config/systemd-system/ly@tty2.service.d/wait-home.conf
[Unit]
# /etc/ly/config.ini ist ein Symlink ins separate /home (Repo unter
# /home/<user>/dotfiles). ly@tty2 hat von Haus aus keine Abhängigkeit zum
# /home-Mount und startet beim Boot fast zeitgleich damit — greift es zu früh,
# liest es die tote Symlink-Config und fällt auf die ly-Defaults zurück
# (u. a. session_log im $HOME, falsche Optik/Optionen). RequiresMountsFor zieht
# home.mount als Ordering- UND Requires-Abhängigkeit ein, sodass ly erst nach
# dem Mount startet.
#
# WICHTIG: Diese Drop-in-Datei muss auf der Root-Partition REAL liegen (nicht
# per links.conf ins /home symlinken) — systemd liest Unit-Drop-ins beim frühen
# Manager-Start, da wäre ein /home-Symlink genauso tot wie die Config selbst.
# Deshalb steht sie nicht in links.conf; sie wird als echte Kopie deployt
# (siehe README »Manual system state«).
RequiresMountsFor=/home
